# scripts/systemd/task-manager-backup.timer
[Unit]
Description=Nightly Personal Task Manager SQLite backup

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
RandomizedDelaySec=300
Unit=task-manager-backup.service

[Install]
WantedBy=timers.target
